# quranbot-dashboard.service
[Unit]
Description=QuranBot Web Dashboard
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/QuranBot/dashboard
Environment=PATH=/root/QuranBot/dashboard/.venv/bin
ExecStart=/root/QuranBot/dashboard/.venv/bin/python run_dashboard.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/root/QuranBot/data /root/QuranBot/logs

# Resource limits
MemoryMax=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target 