# dashboard/deploy/quranbot-dashboard.service
[Unit]
Description=QuranBot Dashboard - Islamic Knowledge Web Interface
After=network.target
Documentation=https://github.com/yourusername/QuranBot

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/path/to/quranbot/dashboard
Environment="PATH=/path/to/quranbot/dashboard/venv/bin"
Environment="FLASK_ENV=production"
ExecStart=/path/to/quranbot/dashboard/deploy/start_dashboard.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quranbot-dashboard

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target